# packaging/debian/tas-agent.prerm
#!/bin/sh
set -e
systemctl disable tas-agent-askpass.path || true
systemctl stop tas-agent-askpass.path tas-agent-askpass.service || true
