# packaging/debian/tas-agent.postinst
#!/bin/sh
set -e
systemctl daemon-reload
systemctl enable tas-agent-askpass.path || true
